# etc/cron.d/policy-bash
PATH=/usr/local/bin:/usr/bin:/bin:/PATH/TO/POLICY/BIN/UTILS
# m h  dom mon dow	user	command
* *	* * *		root	policy-unhold
* *	* * *		root	policy-greylist-expire
